# launchd/com.hosaypeng.apptracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hosaypeng.apptracker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/hsp/.apptracker/tracker.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/hsp/.apptracker</string>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/hsp/.apptracker/logs/tracker.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/hsp/.apptracker/logs/tracker.stderr.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
